# qsdk/package/autelan-shell/src/jsock/etc/xinetd.d/jsock_tcp
service jsock
{
	type = UNLISTED
	socket_type = stream
	protocol = tcp
	port = 9527
	user = root
	wait = no
	server = /etc/jsock/jsock.cb
	disable = no

	env = TZ=UTC-8
}
