# qsdk/package/autelan-shell/src/jsock/etc/xinetd.d/jsock_udp
service jsock
{
	type = UNLISTED
	socket_type = dgram
	protocol = udp
	port = 9527
	user = root
	wait = no
	server = /etc/jsock/jsock.cb
	disable = no

	env = TZ=UTC-8
}
